# s6-rootfs/example/rootfs-s6-rc/etc/s6-overlay/s6-rc.d/01-adguardhome/run
#!/command/execlineb -P
# shellcheck shell=bash

/opt/adguardhome/AdGuardHome --no-check-update -c /opt/adguardhome/conf/AdGuardHome.yaml -w /opt/adguardhome/work