# s6-rootfs/example/rootfs-s6-rc/etc/s6-overlay/s6-rc.d/02-caddy/run
#!/command/execlineb -P
# shellcheck shell=sh

/opt/caddy/caddy run --config /opt/caddy/Caddyfile --adapter caddyfile